# script/crontab
0 */4 * * * cd /root/chain-alert/cex-alert/spot/volume && python3 volume_four_hour.py
0 * * * * cd /root/chain-alert/cex-alert/spot/volume && python3 volunm_one_hour.py